# digdag.service
[Unit]
Description=digdag

[Service]
Restart=on-failure
ExecStart=/opt/digdag/bin/start.sh

[Install]
WantedBy=multi-user.target
